# deploy/fabric/templates/webvirtmgr.ini
[program:webvirtmgr]
autorestart = True
autostart = True
command = /usr/bin/python /var/www/webvirtmgr/manage.py run_gunicorn -c /var/www/webvirtmgr/conf/gunicorn.conf.py
directory = /var/www/webvirtmgr
redirect_stderr = True
stdout_logfile = /var/log/supervisor/webvirtmgr.log
user = nginx
