# build/webvirtmgr.ini
[program:webvirtmgr]
command=/bin/sh -c "bash /opt/webvirtmgr/bootstrap.sh && /usr/bin/gunicorn -c /opt/webvirtmgr/conf/gunicorn.conf.py webvirtmgr.wsgi"
directory=/opt/webvirtmgr
autostart=true
autorestart=true
stdout_logfile=/tmp/sup_webvirtmgr.log
stderr_logfile=/tmp/sup_webvirtmgr_err.log
user=webvirtmgr

[program:webvirtmgr-console]
command=/usr/bin/python /opt/webvirtmgr/console/webvirtmgr-console
directory=/opt/webvirtmgr
autostart=true
autorestart=true
stdout_logfile=/tmp/sup_webvirtmgr-console.log
stderr_logfile=/tmp/sup_webvirtmgr-console_err.log
user=webvirtmgr
